# systemd/camouflage-server@.service
[Unit]
Description=Camouflage server service
After=network.target

[Service]
ExecStart=/usr/bin/camouflage -server /etc/camouflage/%i.toml
MemoryAccounting=true
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
